# systemd/jarvis-web.service
[Unit]
Description=Jarvis Web API (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/workstation/jarvis_assistant
EnvironmentFile=%h/workstation/jarvis_assistant/.env
ExecStart=/usr/bin/env uv run jarvis-web
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
